# clwatch.service
[Unit]
Description=Run clwatch

[Service]
Type=simple
User=username
ExecStart=/home/dir/scripts/./clwatch.sh
